# systemd/dns-watch.service
[Unit]
Description=DNS availability and latency check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dns-watch.sh
